# deploy/config/mapr_cluster_on_k8s/pre_setup_script/startup_scripts/container_init_services/web.service
[Unit]
Description=init mapr cluster

[Service]
Type=simple
ExecStart=/tmp/web-init-script
TimeoutStartSec=0

[Install]
WantedBy=default.target
